# ha-codex/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
set -euo pipefail

exec node /opt/ha-codex/server.js
